# states/units/check-ipv6.service
[Unit]
Description=IPv6 health check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-ipv6.sh
StandardOutput=journal
