# states/units/sing-box-tun-react.path
[Unit]
Description=Watch sing-box runtime config for live restarts

[Path]
PathChanged={{ home }}/.config/sing-box-tun/config.json
Unit=sing-box-tun-react.service

[Install]
WantedBy=multi-user.target
