# states/units/tun6to4.service.j2
[Unit]
Description=6to4 IPv6 tunnel (RFC 3056)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c '\
  ip tunnel add tun6to4 mode sit ttl 64 remote any local {{ public_v4 }} && \
  ip link set tun6to4 up && \
  ip -6 addr add {{ prefix6 }}::1/16 dev tun6to4 && \
  ip -6 route add 2000::/3 via ::192.88.99.1 dev tun6to4'
ExecStop=/bin/sh -c '\
  ip -6 route del 2000::/3 via ::192.88.99.1 dev tun6to4 2>/dev/null || true; \
  ip link set tun6to4 down 2>/dev/null || true; \
  ip tunnel del tun6to4 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
