# states/units/vaultwarden-container.container
{#- Quadlet container unit template for Vaultwarden
    Rendered by container_service macro. Context:
      catalog_entry:    dict from service_catalog.yaml[vaultwarden]
      image:            fully-resolved registry/image@sha256:... reference
      expanded_mounts:  list of {host, container, mode} with ~-expansion applied
      user_scope:       False (system scope)
-#}
[Unit]
Description=Vaultwarden password manager (containerized)
Documentation=https://github.com/dani-garcia/vaultwarden
After=network-online.target local-fs.target
Wants=network-online.target
# StartLimit* keys belong in [Unit] per systemd.unit(5).
StartLimitBurst=5
StartLimitIntervalSec=60

[Container]
Image={{ image }}
ContainerName=vaultwarden
PublishPort=127.0.0.1:{{ catalog_entry.port }}:80
{% for m in expanded_mounts -%}
Volume={{ m.host }}:{{ m.container }}:{{ m.mode }}
{% endfor -%}
{% if catalog_entry.env_file -%}
EnvironmentFile={{ catalog_entry.env_file }}
{% endif -%}
Environment=SIGNUPS_ALLOWED=false
Environment=DOMAIN=https://vw.{{ salt['grains.get']('host') }}.local
# Vaultwarden health endpoint
HealthCmd=curl -sf http://127.0.0.1:80 >/dev/null 2>&1
HealthInterval=5s
HealthStartPeriod=15s
HealthRetries=3
HealthTimeout=10s
Notify=healthy

[Service]
Restart=on-failure
RestartSec=5
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
