# states/units/user/bw-sync.timer
[Unit]
Description=Hourly Vaultwarden-to-gopass sync
[Timer]
OnCalendar=hourly
Persistent=true
[Install]
WantedBy=timers.target
