# states/units/user/cache-cleanup.service
[Unit]
Description=User cache cleanup (paru, pip, npm, flatpak, cargo)

[Service]
Type=oneshot
ExecStart=%h/src/cfg/scripts/cache-cleanup.sh
WorkingDirectory=%h/src/cfg
