# states/units/user/cache-cleanup.timer
[Unit]
Description=Weekly user cache cleanup

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
