# states/units/user/music-index.service
[Unit]
Description=Music index update (Annoy + Essentia)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/music-index
