# states/units/user/rkn-domains-fetcher-retry@.service
[Unit]
Description=RKN Domains Fetcher Retry %i
After=rkn-domains-fetcher.service
PartOf=rkn-domains-fetcher.service

[Service]
Type=oneshot
ExecStart=%h/src/cfg/scripts/rkn-domains-integration.sh manual
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=%h/src/cfg

# Quick retry with different source
ExecStartPre=/bin/sleep 60
Environment=RETRY_SOURCE=https://raw.githubusercontent.com/antifilter/filterlist/master/antifilter.list

# Emergency fallback on second failure
ExecStopPost=/bin/sh -c 'if [ $EXIT_CODE -ne 0 ]; then %h/src/cfg/scripts/rkn-domains-integration.sh emergency; fi'

[Install]
WantedBy=multi-user.target