# states/units/user/rkn-domains-fetcher.service
[Unit]
Description=RKN Blocked Domains Fetcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/src/cfg/scripts/rkn-domains-integration.sh manual
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=%h/src/cfg

# Autocorrection and fallback
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rkn-domains-fetcher

# Cleanup on failure
ExecStopPost=/bin/sh -c 'if [ $EXIT_CODE -ne 0 ]; then %h/src/cfg/scripts/rkn-domains-integration.sh emergency; fi'

[Install]
WantedBy=default.target