# states/units/user/salt-drift.service
[Unit]
Description=Salt comprehensive drift detection (files, packages, units, runtime)

[Service]
Type=oneshot
ExecStart=%h/src/cfg/scripts/drift_state.py full --project-dir %h/src/cfg
WorkingDirectory=%h/src/cfg
