# states/units/user/telethon-bridge-react.service
[Unit]
Description=Safely react to Telethon Bridge runtime file changes
After=default.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/telethon-bridge-react
